# tailscale-gateway.service
[Unit]
Description=Tailscale Gateway Container
# Only start once the network stack is fully up and an address is assigned.
# tailscale-gateway-autonet.service guarantees a default route exists first
# (via DHCP, or via auto-static fallback on a DHCP-less network).
After=network-online.target docker.service tailscale-gateway-autonet.service
Wants=network-online.target tailscale-gateway-autonet.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

# The script detects the local subnet and launches the container.
ExecStart=/usr/local/bin/start-tailscale-gateway.sh

# Graceful stop: let Tailscale deregister before the container is removed.
ExecStop=/bin/bash -c 'docker exec tailscale-gateway tailscale logout 2>/dev/null; docker stop tailscale-gateway'

# If the script fails (e.g. no DHCP yet), retry up to 3 times with a 10s gap.
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
