# tailscale-gateway-radarr.service
[Unit]
Description=Radarr movie library manager (tailnet-only)
# Needs Docker, the tailnet up (to bind to tailscale0), and the NVMe mounted.
After=docker.service tailscale-gateway.service mnt-nvme.mount
Wants=tailscale-gateway.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/start-tailscale-radarr.sh
ExecStop=/bin/bash -c 'docker stop tailscale-gateway-radarr 2>/dev/null || true'
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
